# confs/systemd/finance.service.d/env.conf
[Service]
Environment=SECRET_KEY='<SECRET KEY HERE>'
Environment=BOT_TOKEN='<TG BOT TOKEN HERE>'